# loadreport.cron
* * * * * vagrant /vagrant/loadreport.rb >> /vagrant/loadreport.log 2>&1
