# stockreport.cron
* * * * * vagrant /vagrant/stockreport.py >> /vagrant/stockreport.log 2>&1
